# root/etc/s6-overlay/s6-rc.d/init-restic/run
#!/usr/bin/with-contenv bash

# create folders
mkdir -p \
    /config/logs

# permissions
chown -R abc:abc \
    /config
